# config/files/GRMLBASE/etc/systemd/system/getty@tty2.service.d/override.conf
# This file was deployed via grml-live.

[Service]
Type=idle
Restart=always
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty2
WorkingDirectory=~
ExecStart=
ExecStart=-/usr/share/grml-scripts/run-screen

[Unit]
Description=screen on tty2
