# config/files/GRMLBASE/etc/systemd/system/ssh.service.d/keygen.conf
# This file was deployed via grml-live.

[Service]
ExecStartPre=-/usr/bin/ssh-keygen -A
