# debian/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
  update-rc.d grml-autoconfig remove >/dev/null || exit 0
fi

#DEBHELPER#

exit 0
