# config/files/GRMLBASE/etc/systemd/system/getty@tty1.service.d/override.conf
# This file was deployed via grml-live.

[Service]
Type=simple
Restart=always
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYVTDisallocate=no
WorkingDirectory=~
ExecStart=
ExecStart=-/usr/share/grml-scripts/run-welcome

[Unit]
Description=grml-quickconfig on tty1
After=grml-autoconfig.service
