# grml-autoconfig/debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
  # upgrade from grml-autoconfig before 0.23.3
  update-rc.d grml-autoconfig remove >/dev/null || exit 0
fi

#DEBHELPER#

exit 0
