# grml-autoconfig/debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
  update-rc.d grml-autoconfig start 30 2 . >/dev/null
fi

#DEBHELPER#

exit 0
